# src/build/etc/puppet/modules/phpfarm/templates/php.conf.erb
[program:php]
command=/phpfarm/inst/php-5.5.17-dev/sbin/php-fpm -F
autostart=true
autorestart=true
