# src/php/build/modules/php/files/etc/supervisor/conf.d/php.conf
[program:php]
command = /phpfarm/inst/php-5.5.18/sbin/php-fpm -F
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
